# .config/autostart/Conky.desktop
[Desktop Entry]
Exec=~/.config/conky/start.sh
Name=Conky
Type=Application
