# .config/autostart/DiscordFPRPC.desktop
[Desktop Entry]
Exec=~/.config/DiscordFlatpakRPC.sh
Name=DiscordFPRPC
Type=Application
