# .config/autostart/DiscordFPRPC.desktop
[Desktop Entry]
Type=Application
Name=DiscordFPRPC
Exec=bash -c "~/.config/DiscordFlatpakRPC.sh"
StartupNotify=false
Terminal=false
