# .config/autostart/XDG GTK.desktop
[Desktop Entry]
Exec=exec /usr/libexec/xdg-desktop-portal-gtk
Name=XDG GTK
Type=Application
Version=1.0
