# overlord/.config/autostart/Conky.desktop
[Desktop Entry]
Exec=~/.config/conky/Regulus/start.sh
Name=Conky
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
