# config/autostart/Conky.desktop
[Desktop Entry]
Exec=conky -d
Name=Conky
Type=Application
Version=1.0
Hidden=false
Comment=starts conky
RunHook=0

